# aria/packaging/aria.service
# systemd *user* service so Aria runs in the background and listens for the wake
# word, firing alarms even with no terminal open. As a user service it shares
# your login session's PipeWire audio (no system-service audio headaches).
#
# Install/enable (the `aria` CLI wraps these):
#     aria enable          # == systemctl --user enable --now aria.service
#     aria status / logs / stop / disable
[Unit]
Description=Topol — voice-first AI assistant (background daemon)
# Wait for the user audio session AND network readiness so we don't crash-loop at
# login/boot (DNS isn't up the instant WiFi associates). For --user units the
# network-online.target ordering is best-effort, so the daemon ALSO treats a
# connection error as retryable (backs off, never exits) as the real guarantee.
After=default.target sound.target network-online.target pipewire.service pipewire-pulse.service wireplumber.service
Wants=sound.target network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/aria daemon
Restart=on-failure
RestartSec=5
TimeoutStopSec=10
# Aria reads its key from the OS keyring or, when that's unavailable under
# `systemd --user`, the encrypted-at-rest fallback file in the state dir.
# Inherit the graphical session env so a headful browser (food ordering) can show
# on the user's X11/Wayland display. On Cinnamon/X11 the session normally exports
# these to the --user manager already; this makes the dependency explicit. If
# there's still no display, web ordering degrades with a spoken "I need your screen".
PassEnvironment=DISPLAY XAUTHORITY WAYLAND_DISPLAY

[Install]
WantedBy=default.target
